# openclaw-gateway.vps.service
[Unit]
Description=OpenClaw Gateway (v2026.5.12)
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
# We've switched back to a System Service (rather than --user) 
# to bypass the D-Bus 'No medium found' errors in your VPS environment.
ExecStart=/usr/bin/node /usr/lib/node_modules/openclaw/dist/index.js gateway --port 18789 --config /home/aiagent/.openclaw/openclaw.json
WorkingDirectory=/data/workspace
User=aiagent
Group=aiagent
Restart=always
RestartSec=5
RestartPreventExitStatus=78
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group

Environment=HOME=/home/aiagent
Environment=TMPDIR=/tmp
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/aiagent/.local/bin
Environment=OPENCLAW_CONFIG_PATH=/home/aiagent/.openclaw/openclaw.json
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Gateway"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.5.12

[Install]
WantedBy=multi-user.target
